# ebook-summarizer.service
[Unit]
Description=Ebook Summarizer Streamlit App
After=network.target

[Service]
Type=simple
# UPDATE THESE PATHS to match your setup:
WorkingDirectory=%h/EbookSummarizer
ExecStart=%h/EbookSummarizer/venv/bin/streamlit run app.py --server.port 3002 --server.address 0.0.0.0 --server.headless true
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
